# aaa-modem.service
[Unit]
Description=Flask Application
After=network.target

[Service]
WorkingDirectory=/home/wapa/Documents/ModemInterface
ExecStart=/bin/bash -c 'source /home/wapa/Documents/ModemInterface/modemVenv/bin/activate && /home/wapa/Documents/ModemInterface/modemVenv/bin/python /home/wapa/Documents/ModemInterface/ModemApi.py'
Restart=on-failure
RestartSec=10
Environment=PATH=/home/wapa/Documents/ModemInterface/modemVenv/bin
Environment=FLASK_APP=ModemApi.py
StandardOutput=inherit
StandardError=inherit

[Install]
WantedBy=multi-user.target
